# Pinhole Attack/deploy/cookbooks/pinhole/files/etc/xinetd.d/pinhole
service pinhole
{
    socket_type     = stream
    type            = UNLISTED
    user            = pinhole
    port            = 23464
    server          = /home/pinhole/server.sh
    wait            = no
}

